# Web_data_collection/scrap_flightradar/docker-scripts/crontab
* * * * * Rscript /usr/local/src/flight-scrap/docker-scripts/docker_script.r > /proc/1/fd/1 2>&1
